# .config/supervisord/supervisord.tangojmodelica.conf
[supervisord]
nodaemon=true
logfile=/var/tmp/supervisord.log
pidfile=/var/tmp/supervisord.pid

[program:tango-starter]
command=/bin/bash -c '
  /usr/local/bin/tango_register_device Starter/`hostname --short` Starter tango/admin/`hostname --short` &&
  exec -- /usr/bin/Starter `hostname --short`'
autorestart=true
startretries=5
priority=1

[program:tanks-optimal-control]
command=/bin/bash -c '
  /usr/local/bin/tango_register_device TanksOptimalControl/TOC TanksOptimalControl opt/ctrl/1 &&
  python /opt/masthe-tango-ds-tanks-regulation/ds_tanks/tanks_optimal_control.py TOC -v4'
autorestart=true
startretries=5
priority=5
killasgroup=true
stdout_logfile=/var/tmp/toc.log
